# FIles/autofan.service
[Unit]
Description=Fan Speed Controler using autofan.sh

[Service]
Restart=always
RestartSec=10s
ExecStart=/bin/bash /autofan.sh 
