# Files/autofan.service
[Unit]
Description=Fan speed controller using autofan.sh

[Service]
Restart=always
RestartSec=10s
ExecStart=/bin/bash /autofan.sh 
